# scenarios/10-apt-remote-access-trojan/config/wazuh-rules.xml
<!-- Scenario 10: APT Remote Access Trojan Detection Rules for Wazuh -->
<group name="wcace,apt,rat,c2,">

  <rule id="100460" level="10">
    <if_group>process_execution</if_group>
    <match>python3|python</match>
    <regex>parent_process.*libreoffice|parent_process.*openoffice|parent_process.*thunderbird</regex>
    <description>WCACE S10: Suspicious Python process spawned from document handler</description>
    <mitre>
      <id>T1219</id>
      <id>T1059</id>
    </mitre>
    <group>rat,delivery,</group>
  </rule>

  <rule id="100461" level="12">
    <if_group>process_execution|file_create</if_group>
    <match>crontab|.bashrc|@reboot|.local/share</match>
    <regex>sys_update|hidden.*True|\.local.*\.py</regex>
    <description>WCACE S10: RAT installation with persistence mechanism detected</description>
    <mitre>
      <id>T1219</id>
      <id>T1053</id>
    </mitre>
    <group>rat,persistence,</group>
  </rule>

  <rule id="100462" level="12" frequency="8" timeframe="600">
    <if_matched_group>c2_beacon</if_matched_group>
    <same_source_ip />
    <description>WCACE S10: C2 beaconing pattern detected - regular interval callbacks to external host</description>
    <mitre>
      <id>T1071</id>
      <id>T1219</id>
    </mitre>
    <group>rat,c2,beaconing,</group>
  </rule>

  <rule id="100463" level="14">
    <if_group>process_execution|c2_command</if_group>
    <match>rat_execution|c2_command|tasking</match>
    <description>WCACE S10: RAT command execution via C2 channel</description>
    <mitre>
      <id>T1219</id>
      <id>T1059</id>
    </mitre>
    <group>rat,c2,execution,</group>
  </rule>

  <rule id="100464" level="14">
    <if_group>process_execution|file_access</if_group>
    <match>id_rsa|/etc/shadow|credential|Login Data|passwords</match>
    <description>WCACE S10: Credential access via RAT</description>
    <mitre>
      <id>T1219</id>
      <id>T1003</id>
    </mitre>
    <group>rat,credential_access,</group>
  </rule>

  <rule id="100465" level="12">
    <if_group>process_execution|file_create</if_group>
    <match>tar czf|data_pkg|staging|archive</match>
    <regex>rat_harvester|\.tar\.gz</regex>
    <description>WCACE S10: Data staging for exfiltration via RAT</description>
    <mitre>
      <id>T1074</id>
    </mitre>
    <group>rat,staging,</group>
  </rule>

  <rule id="100466" level="14">
    <if_group>c2_exfiltration|data_transfer</if_group>
    <match>c2_exfiltration|upload|exfil</match>
    <description>WCACE S10: Data exfiltration over C2 channel</description>
    <mitre>
      <id>T1041</id>
    </mitre>
    <group>rat,exfiltration,</group>
  </rule>

  <rule id="100467" level="12">
    <if_group>process_execution|rat_module</if_group>
    <match>keylogger|screenshot|pynput|screen capture</match>
    <description>WCACE S10: Keylogger or screen capture activity detected</description>
    <mitre>
      <id>T1056</id>
      <id>T1113</id>
    </mitre>
    <group>rat,collection,keylogger,</group>
  </rule>

  <rule id="100468" level="8">
    <if_group>process_execution|file_delete</if_group>
    <match>rm -rf|.keys.log|.scr_|cleanup</match>
    <description>WCACE S10: RAT cleanup/anti-forensics activity</description>
    <mitre>
      <id>T1070</id>
    </mitre>
    <group>rat,anti_forensics,</group>
  </rule>

  <rule id="100469" level="15" frequency="4" timeframe="3600">
    <if_matched_group>rat</if_matched_group>
    <description>WCACE S10: APT RAT campaign correlated - Multiple RAT indicators detected</description>
    <mitre>
      <id>T1219</id>
      <id>T1071</id>
    </mitre>
    <group>rat,apt,correlation,critical,</group>
  </rule>

</group>
